# images/debian_jessie/supervisor/nc.conf
[program:nc_l_9999]
command=/bin/nc -l 9999
priority=1
